# rtl_433.conf
# file for supervisord /etc/supervisor/conf.d/rtl_433.conf
[program:rtl_433]
command=/home/pi/rtl_433/build/src/rtl_433 -c /home/pi/rtl_433/rtl_433.conf
autostart=yes
autorestart=yes
startretries=100
stderr_logfile=/var/log/rtl_433/rtl_433.err.log
